# target/default/includes.chroot/etc/systemd/system/disable-drrs-conditionally.service
[Unit]
Description=Disable DRRS conditionally
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/disable-drrs-conditionally

[Install]
WantedBy=multi-user.target
