# custom-packages/digabi-firewall/etc/network/if-pre-up.d/iptablesload
#!/bin/sh

if [ -f "/etc/iptables.rules" ]
then
	iptables-restore < /etc/iptables.rules
fi

exit 0
