# config/includes.chroot/etc/network/if-pre-up.d/ip6tablesload
#!/bin/sh

if [ -f "/etc/ip6tables.rules" ]
then
	ip6tables-restore < /etc/ip6tables.rules
fi

exit 0
